# systemd/iii-inference.service
[Unit]
Description=iii Engine with Inference Worker (Python)
Documentation=https://iii.dev/docs/
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/app/quickstart
Environment=III_HOST=0.0.0.0
ExecStart=/usr/local/bin/iii engine start --config /opt/app/config-inference.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
MemoryMax=8G

[Install]
WantedBy=multi-user.target
